# arch/usbwakeup.service
[Unit]
Description=Enable wake ups by usb

[Service]
ExecStart=/home/oeegor/dev/dotfiles/arch/usbwakeup.sh

[Install]
WantedBy=multi-user.target
